# schedule/crontab
* * * * * root /usr/local/bin/php -q -f /var/www/test-task/artisan schedule:run — no-ansi >> /var/log/cron.log 2>&1